# build_files/systemd/bootc-fetch-apply-updates.service.d/override.conf
[Service]
# Notify Kured that a reboot is required after applying updates
ExecStartPost=touch /var/run/reboot-required

