# build_files/systemd/bootc-fetch-apply-updates.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/bootc upgrade --quiet
# Notify Kured that a reboot is required after applying updates
ExecStartPost=touch /var/run/reboot-required
