# build_files/systemd/ublue-nvctk-cdi.service.d/override.conf
[Unit]
ConditionKernelCommandLine=!nomodeset
ConditionPathExists=/proc/driver/nvidia/version