# build_files/systemd/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
# reset existing scheduling directives
OnBootSec=
OnUnitInactiveSec=
OnCalendar=

# your new schedule:
OnCalendar=Tue 03:00

# keep the randomized 2h delay
RandomizedDelaySec=2h
